# lab4/crontab
*/5 * * * thu ~/lab4/m1
*/5 * * * thu ~/lab4/m1
